# thumbtack-agent/systemd/thumbtack-agent.service
[Unit]
Description=Thumbtack email-harvester agent (Playwright, box-only)
After=network-online.target
Wants=network-online.target
# Cap restart rate so a crash loop can't hammer Thumbtack (these belong in [Unit]).
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
# A GNOME/RDP desktop login injects WAYLAND_DISPLAY/DISPLAY/XAUTHORITY into the
# systemd user manager; Chrome then attaches to the (dead-after-logout) Wayland
# compositor instead of xvfb-run's virtual display and hangs at launch. Scrub them.
UnsetEnvironment=WAYLAND_DISPLAY DISPLAY XAUTHORITY GNOME_SETUP_DISPLAY
# Adjust these paths at deploy (Task 12). node lives under nvm; xvfb-run gives the
# headful Chrome a virtual display on the headless box.
WorkingDirectory=/home/drbartender/projects/os/thumbtack-agent
EnvironmentFile=/home/drbartender/projects/os/thumbtack-agent/.env
ExecStart=/usr/bin/xvfb-run -a /home/drbartender/.nvm/versions/node/v24.16.0/bin/node src/index.js
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target
